# dev.conf
[program:dev]
directory=/root
command=/root/dev/bin/gunicorn stable:app -b localhost:8090
autostart=true
autorestart=true
stderr_logfile=/root/log/dev.err.log
stdout_logfile=/root/log/dev.out.log
