# swagger.conf
[program:swagger]
directory=/root/swagger-ui/
command=python3 -m http.server 5000
autostart=true
autorestart=true
stderr_logfile=/root/log/swagger.err.log
stdout_logfile=/root/log/swagger.out.log
